# systemd/content-stories.timer
[Unit]
Description=Generate STORIES films daily

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
