# systemd/content-videos.timer
[Unit]
Description=Refresh topics + generate SIGNAL videos hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
