# systemd/job-platforms.timer
[Unit]
Description=Discover best job-search platforms per country daily at 04:00

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
